# docker/prod/root/etc/s6-overlay/s6-rc.d/dockerd/run
#!/command/with-contenv bash

exec /usr/bin/dockerd-entrypoint.sh
