# docker/prod/root/etc/s6-overlay/s6-rc.d/jenkins/run
#!/command/with-contenv bash

exec su jenkins -c "bash /usr/local/bin/jenkins.sh"
